# deploy/bpdrive.service.example
[Unit]
Description=度盘
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bpdrive
Environment=BPDRIVE_ADDR=:18088
Environment=BPDRIVE_DATA=/opt/bpdrive/data
# Optional: set credentials here instead of data/config.json.
# Environment=BPDRIVE_BAIDU_APP_KEY=your_baidu_app_key
# Environment=BPDRIVE_BAIDU_SECRET_KEY=your_baidu_secret_key
# Environment=BPDRIVE_BAIDU_REDIRECT_URI=oob
ExecStart=/opt/bpdrive/bpdrive
Restart=always
RestartSec=3
User=www
Group=www

[Install]
WantedBy=multi-user.target
